# cases/deploy/tmp/global/node/192.168.120.122_16790/node-16790.conf
[program:node-16790]
command=/home/platon/trantor_test/node-16790/platon --identity platon --datadir /home/platon/trantor_test/node-16790/data --port 16790 --db.nogc --nodekey /home/platon/trantor_test/node-16790/data/nodekey --cbft.blskey /home/platon/trantor_test/node-16790/data/blskey --config /home/platon/trantor_test/node-16790/config.json --syncmode 'full' --debug --verbosity 5 --rpc --rpcaddr 0.0.0.0 --rpcport 6790 --rpcapi platon,debug,personal,admin,net,web3,txpool --txpool.nolocals
numprocs=1
autostart=false
startsecs=1
startretries=3
autorestart=unexpected
exitcode=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=20
stdout_logfile=/home/platon/trantor_test/node-16790/log/platon.log
