# cases/deploy/tmp/global/node/192.168.120.123_16789/node-16789.conf
[program:node-16789]
command=/home/platon/trantor_test/node-16789/platon --identity platon --datadir /home/platon/trantor_test/node-16789/data --port 16789 --db.nogc --nodekey /home/platon/trantor_test/node-16789/data/nodekey --cbft.blskey /home/platon/trantor_test/node-16789/data/blskey --config /home/platon/trantor_test/node-16789/config.json --syncmode 'full' --debug --verbosity 5 --rpc --rpcaddr 0.0.0.0 --rpcport 6789 --rpcapi platon,debug,personal,admin,net,web3,txpool --txpool.nolocals
numprocs=1
autostart=false
startsecs=1
startretries=3
autorestart=unexpected
exitcode=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=20
stdout_logfile=/home/platon/trantor_test/node-16789/log/platon.log
